# deploy/systemd/gaokao-staging-frontend.service.example
# 安装: sudo cp gaokao-staging-frontend.service.example /etc/systemd/system/gaokao-staging-frontend.service
# 依赖: 已在 /app/mega/frontend 执行过 npm ci && npm run build
# 然后: sudo systemctl daemon-reload && sudo systemctl enable --now gaokao-staging-frontend
#
# Next 通过环境变量 PORT 监听 3100，与 deploy/nginx/mega 中 proxy_pass 一致。

[Unit]
Description=Gaokao Next.js (staging / mega.theyuanxi.cn)
After=network.target gaokao-staging-backend.service

[Service]
Type=simple
User=root
WorkingDirectory=/app/mega/frontend
Environment=NODE_ENV=production
Environment=PORT=3100
ExecStart=/usr/bin/npm run start
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
